# ubuntu-jeos/rc.local
#!/bin/sh -e

bash /usr/local/bin/ip-info
exit 0
